# wwww.conf
[program:wwww]
command=python3 /home/curator2/wwww/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/wwww/test.err.log
stdout_logfile=/home/curator2/wwww/test.out.log